# docker/services/postgresql/supervisord.conf
[supervisord]
nodaemon=true
user=root
loglevel=info
pidfile=/var/run/supervisord.pid

[program:postgresql]
command=postgres -c shared_preload_libraries=pg_stat_statements -c pg_stat_statements.track=all -c max_connections=100 -c shared_buffers=256MB -c effective_cache_size=1GB -c maintenance_work_mem=64MB -c checkpoint_completion_target=0.9 -c wal_buffers=16MB -c default_statistics_target=100 -c random_page_cost=1.1 -c effective_io_concurrency=200 -c listen_addresses=* -c unix_socket_directories=/var/run/postgresql
user=postgres
environment=POSTGRES_DB=%(ENV_POSTGRES_DB)s,POSTGRES_USER=%(ENV_POSTGRES_USER)s,POSTGRES_PASSWORD=%(ENV_POSTGRES_PASSWORD)s,PGDATA=%(ENV_PGDATA)s
priority=1
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=10

[program:pgbouncer]
command=pgbouncer /etc/pgbouncer/pgbouncer.ini
user=postgres
priority=2
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=5
depends_on=postgresql